# src/conf/deploy/resources/service/openhab-powermeter.service
[Unit]
Description=openhab-powermeter

[Service]
User=root
Group=root
Type=simple
Restart=always
RestartSec=60
ExecStart=/opt/openhab/prod/conf/scripts/POWERMETER/powermeter-watch.sh

StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=openhab-powermeter

[Install]
WantedBy=multi-user.target
