# src/conf/deploy/resources/service/openhab-powermeter.service
[Unit]
Description=openhab-powermeter

[Service]
User=root
Group=root
Type=simple
Restart=always
RestartSec=60
ExecStartPre=/bin/sleep 3
ExecStartPre=/usr/bin/pkill powermeter
ExecStartPre=/bin/sleep 30
ExecStart=/opt/openhab/conf/scripts/POWERMETER/powermeter.sh

StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=openhab-powermeter

[Install]
WantedBy=multi-user.target
